# usuario-icpc-service/icpc-clean-homed.service
[Unit]
Description=Maratona Clean Home Daemon

[Service]
User=root
ExecStart=/usr/sbin/icpc-clean-homed

[Install]
WantedBy=multi-user.target
